# src/ColorDesktop.Launcher/Resource/linux.service
[Unit]
Description=colordesktop

[Service]
ExecStart=ColorDesktop.Launcher

[Install]
WantedBy=multi-user.target